# templates/observability.service.j2
[Unit]
Description=PDS Observability Stack (Grafana + Loki + Promtail)
After=docker.service pds.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ pds_data_dir }}
ExecStart=/usr/bin/docker compose --file {{ pds_data_dir }}/docker-compose.observability.yml up --detach
ExecStop=/usr/bin/docker compose --file {{ pds_data_dir }}/docker-compose.observability.yml down

[Install]
WantedBy=default.target
