# templates/pds-stats.service.j2
[Unit]
Description=PDS Stats Collector

[Service]
Type=oneshot
ExecStart={{ pds_data_dir }}/pds-stats.sh
StandardOutput=append:/var/log/pds-stats.log
