# files/pds-stats.timer
[Unit]
Description=Run PDS Stats Collector every minute

[Timer]
OnCalendar=*:*:00
AccuracySec=5s

[Install]
WantedBy=timers.target
